# systemd/infra-maintenance.service
[Unit]
Description=Infra Demo Maintenance Task

[Service]
Type=oneshot
ExecStart=/home/aaliya/linux-infra-intern-assignment/scripts/maintenance.sh

