# systemd/infra-maintenance.service
[Unit]
Description=Infra Demo Maintenance Task

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/bin/env bash /opt/infra-demo/maintenance.sh
